# cookbooks/redis/templates/default/redis.upstart.conf.erb
#!upstart

description "Redis Server"

env USER=<%= node[:redis][:user] %>

start on runlevel [23]
stop on shutdown

exec sudo -u $USER /usr/local/bin/redis-server /etc/redis/redis.conf

respawn